# Idrivelib/cronsetup/gentoo/gte-1.0/idrivecron
#!/sbin/openrc-run
# @IDrive Inc
# Created By : Sabin Cheruvattil
# Modified By: Senthil Pandian

command="perl __LAUNCHPATH__"
pidfile=/var/run/__APP__cron.pid
command_background=true
